# systemd/wfweb.service
[Unit]
Description=wfweb — browser-based remote control for amateur radio
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=pi
ExecStart=/usr/local/bin/wfweb
[Install]
WantedBy=multi-user.target
